# deploy/systemd/thanda-store-xero-contact-access.timer
[Unit]
Description=Run daily Thanda Store Xero contact-access reconciliation

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
Persistent=true
AccuracySec=1min
Unit=thanda-store-xero-contact-access.service

[Install]
WantedBy=timers.target
