# deploy/systemd/thanda-store-renogy-sync.service
[Unit]
Description=Sync Thanda Store Renogy catalogue and warehouse stock
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/thanda-store/thanda-store
EnvironmentFile=/etc/thanda-store-supplier.env
ExecStart=/usr/bin/npm run sync:renogy
