# deploy/systemd/thanda-store-xero-sales-history.timer
[Unit]
Description=Run Thanda Store Xero sales-history sync every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
AccuracySec=1min
Unit=thanda-store-xero-sales-history.service

[Install]
WantedBy=timers.target
