# deploy/systemd/thanda-store-xero-sales-history.timer
[Unit]
Description=Run daily Thanda Store Xero sales-history reconciliation

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
AccuracySec=1min
Unit=thanda-store-xero-sales-history.service

[Install]
WantedBy=timers.target
